# deploy/systemd/callroo-dashboard.service
[Unit]
Description=Call Roo Dashboard
After=network-online.target callroo-printer.service
Wants=network-online.target

[Service]
Type=simple
User=callroo
Group=callroo
WorkingDirectory=/opt/callroo-printer
EnvironmentFile=-/etc/default/callroo-printer
ExecStart=/opt/callroo-printer/.venv/bin/python -m callroo_printer --config /opt/callroo-printer/config.json --dashboard --dashboard-host 127.0.0.1 --dashboard-port 3001 --log-level INFO
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=25

[Install]
WantedBy=multi-user.target
